# debian/ipmctl-monitor.init
#!/usr/bin/env /lib/init/init-d-script
### BEGIN INIT INFO
# Provides:		ipmctl-monitor
# Required-Start:	$syslog $time $remote_fs
# Required-Stop:	$syslog $time $remote_fs
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
# Short-Description: Intel Persistent Memory monitoring
# Description: Intel Persistent Memory monitoring service
### END INIT INFO
DAEMON=/usr/sbin/ipmctl-monitor
DAEMON_ARGS=-d
